# bin/supervisor_queue.ini
[program:omg_queue]
command=./omg_queue.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/omg_queue.err.log
stdout_logfile=/var/log/omg_queue.out.log
