# bin/supervisor_queue.conf
[program:omgqueue]
command=[project_path]/bin/omg_queue.sh
autostart=true
autorestart=true
stderr_logfile=/var/log/omg_queue.err.log
stdout_logfile=/var/log/omg_queue.out.log
